# track_my_ip.timer
[Unit]
Description=Check battery

[Timer]
OnBootSec=300
OnUnitActiveSec=300

[Install]
WantedBy=timer.target
